# gitlab-workhorse.service
[Unit]
Description=GitLab Workhorse
After=gitlab-ce-puma.service
PartOf=gitlab-ce.target

[Service]
Type=simple
User=git
Group=git
WorkingDirectory=/var/lib/gitlab
ExecStart=/usr/sbin/gitlab-workhorse -listenUmask 0 -listenNetwork unix -listenAddr /run/gitlab/gitlab-workhorse.socket -authBackend http://127.0.0.1:3003 -authSocket /run/gitlab/puma.sock
Restart=on-failure

[Install]
WantedBy=multi-user.target
